# kubernetes/builder/refonte/supervisord.conf
[supervisord]
nodaemon = false

[program:php5-fpm]
command=php-fpm
stdout_events_enabled=true
stderr_events_enabled=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_events_enabled=true
stderr_events_enabled=true